# etc/supervisor/run_redis.conf
[program:run_redis]
directory=/root/redis-stable
command=bash -c "redis-server"
autostart=true
autorestart=true
redirect_stderr=true