# supervisor/my_qcluster.conf
[program:django-q-my_qcluster]
command = /home/ubuntu/.virtualenvs/django_telegram/bin/python /home/ubuntu/parsing_news_and_send_telegram/manage.py my_qcluster
stopasgroup = true