# supervisor/run_django.conf
[program:run_django]
command=/home/ubuntu/run
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
