# supervisor/run_commands_with_server.conf
[program:run_commands_with_server]
command=/home/ubuntu/parsing_news_and_send_telegram/scripts/run.sh
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
