# etc/supervisor/develop/run_commands_with_server.conf
[program:run_commands_with_server]
command=scripts/run.sh
directory=/home/vagrant/parsing_news_and_send_telegram/
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
